# cvmfs-autofs/supervisor_autofs.conf
[program:autofs]
command=/usr/sbin/automount --foreground --dont-check-daemon 
autorestart=true
