# cvmfs-autofs-csi/supervisor_csi.conf
[program:csi]
command=/usr/sbin/csi-cvmfsplugin --nodeid=$(NODE_ID) --endpoint=$(CSI_ENDPOINT) --v=5 --drivername=csi-cvmfsplugin 
autorestart=true
